# app/supervisor.ini
[program:app-uwsgi]
user=root
command=uwsgi --ini /app/uwsgi.ini
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:app-maketopics]
user=root
command=maketopics
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:app-storeauthtable]
user=root
command=storeauthtable
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:app-readfeed]
user=root
command=readfeed -d --sleeptime 600
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:crond]
user=root
command=crond -f
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true


